# infra/oracle/stormwatch.service
[Unit]
Description=StormWatch API + worker (FastAPI/uvicorn)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/storm-watch/backend
EnvironmentFile=/home/ubuntu/storm-watch/backend/.env
# IMPORTANTE: um único worker. O loop de satélite roda dentro do processo;
# mais de um worker enviaria notificações duplicadas.
ExecStart=/home/ubuntu/storm-watch/backend/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000 --workers 1
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
